# docker/shadowsocks.conf
[program:shadowsocks] 
command=python /opt/shadowsocks/shadowsocks/server.py -c /opt/shadowsocks/shadowsocks/config.json
autorestart=true 
user=root 